# virtualizations/rc.local
#!/bin/bash
sh $HOME/start-vms-nogui.sh
